# deploy/adaptive-heating-mvp.service
[Unit]
Description=Adaptive Heating MVP controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jack
WorkingDirectory=/home/jack/adaptive-heating-mvp
Environment=RUST_LOG=adaptive_heating_mvp=info
Environment=PATH=/home/jack/.cargo/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/etc/adaptive-heating-mvp.env
LoadCredential=influx_token:/etc/adaptive-heating-mvp/influx.token
ExecStart=/home/jack/adaptive-heating-mvp/target/release/adaptive-heating-mvp --config /home/jack/adaptive-heating-mvp/model/adaptive-heating-mvp.toml run
ExecStop=/home/jack/adaptive-heating-mvp/target/release/adaptive-heating-mvp --config /home/jack/adaptive-heating-mvp/model/adaptive-heating-mvp.toml restore-baseline
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
