# scripts/dhw-auto-trigger.service
[Unit]
Description=DHW auto-trigger — forces eBUS DHW charge on sustained hot water draw
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/dhw-auto-trigger.sh
Restart=always
RestartSec=10

# Run unprivileged — only needs network access
User=jack
Group=jack

# Logging goes to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dhw-auto-trigger

[Install]
WantedBy=multi-user.target
